# etc/udev/rules.d/75-mac-spoof.rules
ACTION=="add", SUBSYSTEM=="net", ATTR{address}=="HWWWWWWWWWWWWW", RUN+="/usr/bin/ip link set dev %k address notHWWWWWWWWWWWWW"
ACTION=="add", SUBSYSTEM=="net", ATTR{address}=="HWWWWWWWWWWWWW", RUN+="/usr/bin/ip link set dev %k address notHWWWWWWWWWWWWW"
