# home/.config/autostart/tilda.desktop
[Desktop Entry]
Name=Tilda
Comment=TildaTerm
Exec=tilda
Type=Application
Terminal=false
Hidden=false

